# examples/systemd/phase05-status-job.service.example
[Unit]
Description=Phase 5 Status Job

[Service]
Type=oneshot
User=<vm-user>
WorkingDirectory=/home/<vm-user>
ExecStart=/usr/local/bin/phase05-status-job.sh systemd-timer
