# infra/virtualbox/phase04/phase04-heartbeat.service
[Unit]
Description=Phase 4 Heartbeat Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/phase04-heartbeat.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
